# Vagrant/Flask_Postgres/cookbooks/flask/files/default/service.conf
start on started sshd
stop on runlevel [!2345]

exec /usr/bin/python2.6 /home/vagrant/greedy.py
respawn